# files/laptop_power_management/udev/99-50-dc-power.rules
# AC UNPLUGGED
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/dc_powermode.sh"
