# files/laptop_power_management/udev/99-60-dc-critical-power.rules
# Low Battery
SUBSYSTEM=="power_supply", ATTR{capacity}=="[0-1][0-5]", RUN+="/usr/local/bin/dc_critical_powermode.sh"
